# script/rc.local
find /etc/rc.d/rc.local.d/ -type f -exec /bin/sh {} \;
